# etc/udev/rules.d/98-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/bin/systemctl --no-block start hotplug-monitor.service"
