# etc/systemd/system/hotplug-monitor.service
[Install]
WantedBy=multi-user.target

[Unit]
 Description=Hotplug monitor

[Service]
 Type=simple
 ExecStart=/home/lukas/Scripts/monitor.sh
